# src/crontab
* * * * * bash /home/viveorange/lanzar_cron.sh >> /var/log/cron.log 2>&1
